# airootfs/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Unit]
Requires=pacman-init.service
After=pacman-init.service

[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -f -- \\u' --noclear --autologin root - $TERM
